# telemetrics-client-1.2.0/src/data/journal-probe.service.in
[Unit]
Description=Telemetrics Systemd Journal Probe
Requires=telemd.socket
After=telemd.socket

[Service]
ExecStart=@bindir@/journalprobe
User=telemetry

[Install]
WantedBy=multi-user.target
