# telemetrics-client-1.2.0/src/data/telemd.socket.in
[Unit]
Description=Telemetrics Daemon

[Socket]
ListenStream=@SOCKETDIR@/telem-0

[Install]
WantedBy=sockets.target
