# telemetrics-client-1.2.0/src/data/telemd.service.in
[Unit]
Description=Telemetrics Daemon

[Service]
ExecStart=@bindir@/telemd
User=telemetry

[Install]
Also=telemd.socket telemd.path
WantedBy=multi-user.target
