# telemetrics-client-1.2.0/src/data/telemd.path.in
[Unit]
Description=Telemetry Daemon Spool

[Path]
DirectoryNotEmpty=@localstatedir@/spool/telemetry

[Install]
WantedBy=multi-user.target
